# devopsfetch.service

# devopsfetch.service

[Unit]
Description=DevOpsFetch Service
After=network.target

[Service]
ExecStart=/home/ubuntu/HNG-Task-5/devopsfetch.sh -d all
User=ubuntu  
WorkingDirectory=/home/ubuntu/HNG-Task-5
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=devopsfetch

[Install]
WantedBy=multi-user.target

